# dynv6-yourdomain.timer
[Unit]
Description=Timer for script-dynv6_yourdomain.service

[Timer]
OnBootSec=2m
OnUnitActiveSec=2m

[Install]
WantedBy=multi-user.target